# arch-linux/auto-update/auto-flatpak-update.service
[Unit]
Description=Automatic Flatpak Updates

[Service]
Type=oneshot
ExecStart=/usr/local/bin/auto-flatpak-update.sh
