# arch-linux/auto-update/auto-pacman-update.service
[Unit]
Description=Automatic Pacman Updates

[Service]
Type=oneshot
ExecStart=/usr/local/bin/auto-pacman-update.sh
